# RaspberryPi/controller.service
[Unit]
Description=Greenhouse Controller

[Service]
ExecStart=/root/iot/gh-sensors/controller.py

[Install]
WantedBy=multi-user.target
