# DellEdgeGateway/readsensors.service
[Unit]
Description=Greenhouse Readsensors

[Service]
ExecStart=/opt/greenhouse/readsensors.py

[Install]
WantedBy=multi-user.target
